# templates/node_exporter.service.j2
[Unit]
Description=Node Exporter

[Service]
ExecStart={{ node_exporter_dir }}/node_exporter

[Install]
WantedBy=default.target
